# crontab
#minute	hour	day	month	weekday	command
*/3	*	*	*	*	/usr/bin/healer.sh
